# deploy/weather.service
[Unit]
Description=Weather Server daemon
After=network.target
Documentation=https://github.com/rmrobinson/weather-server

[Service]
Type=simple
User=weather
Group=weather
WorkingDirectory=/etc/weather

# Secrets and runtime config — create /etc/weather/environment and populate:
#   INFLUX_TOKEN=<your-influxdb-api-token>
#   WEATHER_PSK=<your-psk-or-leave-absent-to-disable-auth>
EnvironmentFile=/etc/weather/environment

ExecStart=/usr/local/bin/weather --config /etc/weather/config.yaml

Restart=on-failure
RestartSec=5

# Harden the service
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=

[Install]
WantedBy=multi-user.target
